# slackware/packages/jdk/jdk-8u172/profile.d/jdk.sh
version https://git-lfs.github.com/spec/v1
oid sha256:7fdb1d9564e085fa5a44de631c0546384fea1aea05428caf85241d846210daed
size 146
