# slackware/packages/jdk/jdk-8u172/profile.d/jdk.csh
version https://git-lfs.github.com/spec/v1
oid sha256:494ee8a2a932560738e5c6a4cf740d523c730282fcf3354397b34f6846c47a1f
size 143
